# config/crontab
# TEST: Run 1 minute from now (minute 32 15)
32 15 * * * /app/scripts/cron-ocr-job.sh

# Hourly OCR processing (disabled during test)
# 0 * * * * /app/scripts/cron-ocr-job.sh

# Ensure newline at end
