# config/crontab
# Set PATH and SHELL for cron jobs
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
SHELL=/bin/bash

# Daily OCR processing - runs at midnight
0 0 * * * /app/scripts/cron-ocr-job.sh

# Ensure newline at end
