# config/crontab
# Set PATH and SHELL for cron jobs
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
SHELL=/bin/bash

# OCR processing - runs every hour
# - Skips files modified <60s ago (prevents processing mid-sync)
# - Skips files uploaded in last 8 hours (prevents sync conflicts)
# - Skips already-processed files (fast hash comparison)
0 * * * * /app/scripts/cron-ocr-job.sh

# Full OCR run at 3:30am - processes ALL files regardless of upload time
# Low conflict risk since user is likely asleep
# Runs at :30 to avoid overlap with hourly job at :00
30 3 * * * SKIP_RECENT_CHECK=true /app/scripts/cron-ocr-job.sh

# Ensure newline at end
