# config/crontab
# Set PATH and SHELL for cron jobs
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
SHELL=/bin/bash

# OCR processing - runs every minute
# Safe to run frequently:
# - Skips files modified <10s ago (prevents processing mid-sync)
# - Skips already-processed files (fast hash comparison)
# - Idle overhead: ~1MiB RAM, 0% CPU
* * * * * /app/scripts/cron-ocr-job.sh

# Ensure newline at end
