# config/crontab
# Set PATH and SHELL for cron jobs
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
SHELL=/bin/bash

# OCR processing - runs every 10 minutes
# Safe to run frequently: gracefully skips files that don't need processing
# Database updates happen while sync server runs (no restart needed)
*/10 * * * * /app/scripts/cron-ocr-job.sh

# Ensure newline at end
